# johntrading.service
[Unit]
Description=John Trading System - Automated 24/7 Real-Time Trading
Documentation=https://github.com/Manoj-1212/johnTrading
After=network-online.target
Wants=network-online.target

# Restart policy - restarts if fails too many times
StartLimitIntervalSec=3600
StartLimitBurst=5

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/johntrading

# Load API keys and configuration from .env file
EnvironmentFile=/home/ubuntu/johntrading/.env

# Python environment setup
Environment="PATH=/home/ubuntu/johntrading/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="PYTHONUNBUFFERED=1"
Environment="PYTHONPATH=/home/ubuntu/johntrading"
Environment="TRADING_MODE=paper"

# Main process - Automated 24/7 trading (exits at market close, restarts via cron/systemd)
ExecStart=/home/ubuntu/johntrading/.venv/bin/python /home/ubuntu/johntrading/run_phase9_production.py

# Restart behavior
Restart=always
RestartSec=30

# Resource limits
MemoryMax=512M
CPUQuota=50%

# Logging to systemd journal
StandardOutput=journal
StandardError=journal
SyslogIdentifier=johntrading

# Process management
KillMode=mixed
TimeoutStopSec=30

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/home/ubuntu/johntrading

[Install]
WantedBy=multi-user.target
